# Target/S99UltiFi
#!/bin/sh /etc/rc.common

START=99
start() {
	/UltiFi/PrinterManager.sh > /dev/null 2> /dev/null &
}
